# docker-config/template/elasticsearch/supervisord.conf
[supervisord]
nodaemon=true

[program:apply_entrypoint]
command=bash /docker-entrypoint.sh
autostart=true
autorestart=true

[program:run]
command=elasticsearch
autostart=true
autorestart=true

